# reboot.service
[Unit]
Description=reboot job

[Service]
Type=oneshot
ExecStart=/sbin/reboot
